# devil-advocate-bot.service
[Unit]
Description=Devil's Advocate Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/devil-advocate-bot
EnvironmentFile=/opt/devil-advocate-bot/.env
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
